# netatalk/rootfs/etc/services.d/afpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start afpd service
# ==============================================================================

bashio::log.info "Starting afpd ..."
exec afpd -d
