# rsync/rootfs/etc/services.d/rsyncd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start smbd service
# ==============================================================================

bashio::log.info "Starting rsyncd service ..."
exec rsync --config=/etc/rsyncd.conf --daemon --no-detach
