# rsync/rootfs/etc/services.d/rsyncd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start rsyncd service
# ==============================================================================

bashio::log.info "Starting rsyncd service ..."
exec rsync --daemon --no-detach --config=/config/rsyncd.conf
