# traefik/rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start traefik service
# ==============================================================================

declare config

config=$(bashio::config.true 'custom_config' && echo "/config/traefik.yaml" || echo "/etc/traefik/traefik.yaml")

bashio::log.info "Create additional environment variables ..."
for id in $(bashio::config 'env_vars|keys'); do
  env_var=$(bashio::config "env_vars[${id}]")
  export $env_var
  bashio::log.info "-> export $env_var"
done

if bashio::config.is_empty "env_vars"; then
  bashio::log.info "No additional environment variables found."
fi

bashio::log.info "Starting Traefik ..."
exec traefik --configFile=$config
