# webdav/rootfs/etc/services.d/webdav/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start lighttpd service
# ==============================================================================

bashio::log.info "Starting lighttpd service ..."
if bashio::config.true 'custom_config'; then
  exec lighttpd -D -f /config/lighttpd.conf
else
  exec lighttpd -D -f /etc/lighttpd/lighttpd.conf
fi