# netatalk-conf/rootfs/etc/services.d/avahi-publish-service-1/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
declare hostname

# Read hostname from API or setting default "hassio"
hostname=$(bashio::info.hostname)
if bashio::var.is_empty "name"; then
    bashio::log.warning "Can't read hostname, using default."
    hostname="homeassistant"
fi
bashio::log.info "Starting the AVAHI daemon for ${hostname}..."

exec avahi-publish-service -v -f -s $hostname _afpovertcp._tcp 548
